# conf/supervisord.calamari.ini
[program:carbon-cache]
command=python /usr/bin/carbon-cache --config=/etc/graphite/carbon.conf --logdir=/var/log/carbon/ --pidfile=/var/run/carbon-cache.pid start
user=root

[program:cthulhu]
command=python /usr/bin/cthulhu-manager
user=root