# conf/supervisord.calamari.ini
[program:carbon-cache]
command=python /usr/bin/carbon-cache --config=/etc/graphite/carbon.conf --pidfile=/var/run/carbon-cache.pid start
user=apache

[program:cthulhu]
command=python /usr/bin/cthulhu-manager
user=root