# conf/supervisord.production.conf
[program:carbon-cache]
command=/usr/bin/carbon-cache --config=/etc/graphite/carbon.conf --pidfile=/var/run/carbon-cache.pid start

[program:cthulhu]
command=/usr/bin/cthulhu-manager